# images/postfix/config/saslauthd.conf
[program:saslauthd]
command=/usr/sbin/saslauthd -a pam -n 2 -c -m /var/spool/postfix/var/run/saslauthd
autostart=true
autorestart=true
startsecs=10
startretries=5